# nginx.cron
*/10  *   *   *   *   root    /usr/sbin/logrotate /etc/logrotate.d/nginx
